# wayne/startup/autostart/wayne.service
[Unit]
Description=W.A.Y.N.E - Wireless Artificial Yielding Network Engine
After=network.target sound.target
Wants=network.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/path/to/wayne
ExecStart=/usr/bin/python3 /path/to/wayne/startup/wayne_daemon.py
Restart=always
RestartSec=5
Environment=DISPLAY=:0
Environment=PULSE_RUNTIME_PATH=/run/user/1000/pulse
EnvironmentFile=/path/to/wayne/backend/.env

[Install]
WantedBy=multi-user.target
